# restic/crontab
# This configuration runs "bash /home/pi/dobackup.sh"
# From left to right: every 8 hours on the 0th minute, every day, every month, and every week.
0 8 * * * bash /home/pi/dobackup.sh